# config/autostart/piheater.desktop
[Desktop Entry]
Type=Application
Name=PiHeater
Exec=/usr/local/bin/piheater
