# config/autostart/pinode.desktop
[Desktop Entry]
Type=Application
Name=PiNode
Exec=/usr/local/bin/pinode_tempsensor
